# auto-xorg.service
[Unit]
Description=Auto-Xorg

[Service]
ExecStart=/bin/bash /usr/sbin/auto-xorg.bash
# ExecStart=/bin/bash /usr/sbin/auto-xorg.bash -l -i
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=multi-user.target