# auto-xorg.service
[Unit]
Description=Auto-Xorg

[Service]
ExecStart=/bin/bash /usr/local/bin/auto-xorg
# ExecStart=/bin/bash /usr/local/bin/auto-xorg -l -i
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=multi-user.target